# meta-ai-assistant/recipes-ai/ai-assistant/files/ai-assistant.service
[Unit]
Description=AI Rover Assistant
After=can0.service network-online.target
Wants=network-online.target
Requires=can0.service

[Service]
Type=simple
Environment="GEMINI_API_KEY=your-key-here"
WorkingDirectory=/opt/ai-assistant
ExecStartPre=/bin/sh -c 'if [ ! -d /opt/ai-assistant/venv ]; then python3 -m venv /opt/ai-assistant/venv && /opt/ai-assistant/venv/bin/pip install google-genai opencv-python-headless flask python-can; fi'
ExecStart=/opt/ai-assistant/venv/bin/python3 /opt/ai-assistant/autonomy_web_can.py --channel can0 --camera-interval 0.02 --gemini-interval 1.2
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
